# dist/linux/workload-service.service
[Unit]
Description=Workload Service

[Service]
Type=simple
User=wls
Group=wls
WorkingDirectory=/opt/workload-service/bin/
ExecStart=/opt/workload-service/bin/workload-service startserver
ExecStop=/bin/kill -s SIGTERM $MAINPID
TimeoutStartSec=0
Restart=on-failure
PermissionsStartOnly=true
RuntimeDirectory=workload-service
RuntimeDirectoryMode=0775

[Install]
WantedBy=multi-user.target
