# dist/linux/workload-service.service
[Unit]
Description=Workload Service

[Service]
Type=simple
User=wls
Group=wls
WorkingDirectory=/opt/workload-service/bin/
ExecStart=/opt/workload-service/bin/workload-service startserver
TimeoutSec=300
KillMode=process
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60

[Install]
WantedBy=multi-user.target